# base/univention-network-manager/etc/network/if-pre-up.d/10_dhclientconf
#!/bin/sh
#
# Univention Network
#  Commit dhclient.conf to ensure the expire date is correct
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2004-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only

[ -e /etc/dhcp/dhclient.conf ] || exit 0

. /etc/network/ucs-network-tools
ucs_ignore_interface

ucr commit /etc/dhcp/dhclient.conf

exit 0
