# management/univention-management-console/debian/univention-management-console-server.postinst
#!/bin/sh
#
# Univention Management Console Server
#  postinst script for UMC server
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2002-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only

#DEBHELPER#

# shellcheck source=/dev/null
. /usr/share/univention-lib/all.sh

# set default UCR settings
univention-config-registry set \
	umc/server/debug/level?2 \
	umc/server/autostart?yes \
	umc/server/upload/max?2048 \
	umc/module/debug/level?2 \
	umc/module/timeout?600

create_logfile /var/log/univention/ec2.log "root:adm" 640

# create UMC ACLs for user root
eval "$(ucr shell server/role)"

# append new command patterns to allow list
/usr/sbin/univention-management-console-acls allow --user root --host "*" \
			--command 'ucr/*' \
			--command 'service/*' \
			--command 'quota/*' \
			--command 'top/*' \
			--command 'reboot/*' \
			--command 'services/*' \
			--command 'sysinfo/*' \
			--command 'updater/*' \
			--command 'lib/*'  \
			--command 'join/*'

/usr/sbin/univention-management-console-acls allow --user root --host "*" --flavor "*" \
			--command 'setup/*'

systemctl restart univention-directory-listener

if [ "$1" = "configure" ] && dpkg --compare-versions "$2" lt-nl 12.0.20-1; then
	ucr unset security/packetfilter/package/univention-management-console-server/tcp/6670/all \
		security/packetfilter/package/univention-management-console-server/tcp/6670/all/en
	[ -x "/etc/init.d/univention-firewall" ] && invoke-rc.d univention-firewall restart
fi

# protect secret contents!
touch /usr/share/univention-management-console/oidc/oidc.json
chmod 0600 /usr/share/univention-management-console/oidc/oidc.json

# ucslint needs this
exit 0
