# base/univention-network-manager/etc/network/if-up.d/05_rdate
#!/bin/sh
#
# Univention Network Manager
#  network script: fix time
#  the network link is active
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2008-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only

[ -f /var/run/ntpd.pid ] && exit 0
[ -d /run/systemd/system ] && exit 0

. /etc/network/ucs-network-tools
ucs_ignore_interface

invoke-rc.d rdate start

exit 0
