# base/univention-network-manager/etc/network/if-up.d/30_univention-directory-policy
#!/bin/sh
#
# Univention Network Manager
#  network script: calls univention-directory-policy
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2009-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only

[ -x /etc/init.d/univention-directory-policy ] || exit 0
[ -d /run/systemd/system ] && exit 0

. /etc/network/ucs-network-tools
ucs_ignore_interface

invoke-rc.d univention-directory-policy start >/dev/null 2> /dev/null

exit 0
