# base/univention-network-manager/etc/network/if-up.d/90_dns_update
#!/bin/sh
#
# Univention Network
#  network script: register IP in LDAP
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2004-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only
#

[ -e /var/univention-join/joined ] || exit 0

. /etc/network/ucs-network-tools
ucs_ignore_interface

eval "$(ucr shell server/role ldap/master)"

# Is Primary Directory Node available?
if /bin/netcat -q0 -w1 "$ldap_master" 6670 </dev/null; then

	# If we are the Primary Directory Node a restart of UMC is required
	# PMH: why?
	if [ "$server_role" = "domaincontroller_master" ]; then
		[ -d /run/systemd/system ] ||
		deb-systemd-invoke restart univention-management-console-server
	fi

	if [ -x /usr/sbin/univention-register-network-address ]
	then
	/usr/sbin/univention-register-network-address
	fi
fi

exit 0
