# mail/univention-spamassassin/debian/univention-spamassassin.postrm
#!/bin/sh
#
# Univention Spamassassin
#  postrm script
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2004-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only
# postrm script for univention-spamassassin

#DEBHELPER#

if [ "$1" = "remove" ]; then
	ucr unset mail/antivir/spam
	test -x /etc/init.d/amavis && invoke-rc.d amavis restart
fi

exit 0
