# management/univention-ldap-overlay-memberof/debian/univention-ldap-overlay-memberof.postinst
#!/bin/bash
#
# Univention LDAP Server
#  postinst script for the debian package
#
# Like what you see? Join us!
# https://www.univention.com/about-us/careers/vacancies/
#
# SPDX-FileCopyrightText: 2008-2025 Univention GmbH
# SPDX-License-Identifier: AGPL-3.0-only

#DEBHELPER#

univention-config-registry set ldap/overlay/memberof?true

test -x /etc/init.d/slapd && invoke-rc.d slapd crestart

# new installations
if [ "$1" = "configure" ] && [ -z "$2" ]; then
	# update groups if slapd is running, memberOf is activated and system is joined
	if pidof slapd >/dev/null; then
		# shellcheck source=/dev/null
		. /usr/share/univention-lib/ucr.sh
		test -e /var/univention-join/joined && /usr/share/univention-ldap-overlay-memberof/univention-update-memberof
	fi
fi

exit 0
